# deploy/systemd/robot-video-stack.service
[Unit]
Description=Robot host bridge and RTSP video manager stack
Wants=network-online.target mediamtx.service
After=network-online.target mediamtx.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/robot-video-stack
ExecStart=/bin/bash -lc '${HOSTCOMPUTER_DIR:-/home/rera/hostcomputer}/ros1_bridge/video_stack.sh start'
ExecStop=/bin/bash -lc '${HOSTCOMPUTER_DIR:-/home/rera/hostcomputer}/ros1_bridge/video_stack.sh stop'
TimeoutStartSec=30
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
